# munin-plugins/transmission_xfer.crontab
# m h  dom mon dow   command
*/5 * * * * cp /home/santa/.config/transmission/stats.json /home/santa/.config/transmission/stats_shadow_copy.json && chmod 644 /home/santa/.config/transmission/stats_shadow_copy.json

